# sample_nginx/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
autorestart=true
command=/usr/sbin/nginx -g "daemon off;"

[program:sshd]
autorestart=true
command=/usr/sbin/sshd -D
